# ops/bouwflow-deploy-poller.service
[Unit]
Description=Deploy the tested BouwFlow production release
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=deploy
Group=deploy
Environment="PATH=/opt/node-v22/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
ExecStart=/usr/local/bin/bouwflow-deploy-poller

[Install]
WantedBy=multi-user.target
